# operational_scripts/crontab
0 6 * * * cd ./gis/process_data; /bin/bash -l ./daily_cron.sh >> /var/log/cron/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
